# mansur.service
[Unit]
Description=Log into telegram

[Service]
Type=oneshot
EnvironmentFile=/etc/mansur/vars.env
ExecStart=/bin/mansur

